# app/crontab
#* * * * * /usr/local/bin/python /app/digest.py #Uncomment this line for testing to send emails every minute
0 7,19 * * * /usr/local/bin/python /app/digest.py

